# borg-backup.timer
[Unit]
Description=Hourly Borg Backup

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
