# borg-backup.service
[Unit]
Description=Borg Backup

[Service]
Type=oneshot
ExecStart=/usr/local/etc/backup_by_borg.sh
