# jxy-memories.service
[Unit]
Description=JxyMemories

[Service]
Type=oneshot
ExecStart=/usr/local/etc/jxy-memories.py
