# jxy-memories.timer
[Unit]
Description=Hourly JxyMemories

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
